# mkosi.profiles/gcp/mkosi.extra/etc/systemd/system/reth-sync.service.d/gcp-override.conf
[Service]
Environment=RCLONE_TRANSFERS=5 \
            RCLONE_STREAMS=5
