# mkosi.profiles/playground/mkosi.extra/etc/systemd/system/haproxy.service.d/91-playground.conf
[Unit]
After=render-config.service

[Service]
WorkingDirectory=/etc/haproxy
ExecStartPre=+/usr/lib/playground/generate-haproxy-certs.sh
