# modules/cvm-base/mkosi.extra/etc/systemd/system/cvm-provision.service
[Unit]
Description=CVM Provisioner — mTLS API for /init, /manifest, /reboot, /status
After=network-online.target wait-for-key.service
Wants=network-online.target
Requires=wait-for-key.service

[Service]
Type=simple
ExecStart=/usr/bin/cvm-provisioner \
    --listen 0.0.0.0:8888 \
    --listen-loopback 127.0.0.1:8889 \
    --runtime-dir /run/cvm-provisioner \
    --persistent-mount /persistent \
    --tdx-init-binary /usr/bin/tdx-init \
    --authorized-pubkey-file /etc/searcher_key \
    --mode auto
RuntimeDirectory=cvm-provisioner
RuntimeDirectoryMode=0700
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=minimal.target
