# base/mkosi.skeleton/etc/systemd/system/minimal.target
[Unit]
Description=Minimal System
Requires=basic.target
Conflicts=rescue.service rescue.target emergency.service emergency.target
After=basic.target rescue.service rescue.target emergency.service emergency.target
AllowIsolate=yes

[Install]
WantedBy=default.target