# mkosi.profiles/devtools/mkosi.extra/etc/systemd/system/cvm-provision.service.d/dev-mode.conf
[Service]
# Local QEMU has no TDX; auto falls back to mock for the extend leg while
# keeping the LUKS bring-up real (tdx-init binary is still present).
Environment=CVM_PROVISIONER_MODE=auto
