# mkosi.profiles/playground/mkosi.extra/etc/systemd/system/serial-getty@hvc0.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin bnet --noclear %I $TERM
