# buildernet/mkosi.extra/etc/systemd/system/rbuilder-rebalancer.service.d/90-override.conf
[Service]
EnvironmentFile=%E/rbuilder-rebalancer/rbuilder-rebalancer.env
