# buildernet/mkosi.extra/etc/systemd/system/rbuilder-rebalancer.service.d/90-override.conf
[Unit]
After=render-config.service

[Service]
EnvironmentFile=%E/rbuilder-rebalancer/rbuilder-rebalancer.env
