# modules/l2/_common/mkosi.extra/etc/systemd/system/ptlb-routes-nanny.timer
[Unit]
Description=Ensure presence of local routes required for passthrough load-balancer

[Timer]
OnCalendar=*-*-* *:*:5
Unit=ptlb-routes-nanny.service

[Install]
WantedBy=multi-user.target
