# modules/l2/_common/mkosi.extra/etc/systemd/system/logrotate.timer.d/override.conf
[Timer]
OnCalendar=*-*-* *:0/5:0
RandomizedDelaySec=15
