# mkosi.images/buildernet/mkosi.extra/etc/systemd/system/config-watchdog.timer
[Unit]
Description=Check for new rbuilder and rbuilder-bidding configs

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
