# mkosi.profiles/playground/mkosi.extra/etc/systemd/system/persistent-setup.service.d/playground-override.conf
[Service]
Environment=
ExecStart=
ExecStart=/usr/bin/persistent-setup --device-lun 10 --key-file /etc/disk-encryption-key
TimeoutStartSec=120
