# bob-l2/mkosi.extra/etc/systemd/system/searcher-container.service.d/needs-metadata.conf
[Unit]
After=fetch-metadata.service
Requires=fetch-metadata.service

[Service]
EnvironmentFile=/etc/metadata.env
