# modules/flashbox/flashbox-l2/mkosi.extra/etc/systemd/system/searcher-firewall.service.d/needs-config.conf
[Unit]
After=fetch-config.service
Requires=fetch-config.service

[Service]
EnvironmentFile=/etc/bob/config.env
