# buildernet/mkosi.extra/etc/systemd/system/ssh.socket.d/port40192.conf
[Socket]
ListenStream=
ListenStream=0.0.0.0:40192
