# bob-l1/mkosi.extra/etc/systemd/system/searcher-firewall.service.d/cloud-provision.conf
[Unit]
After=azure-complete-provisioning.service
Wants=azure-complete-provisioning.service