# bob-l1/mkosi.extra/etc/systemd/system/searcher-firewall.service.d/require-azure-provision.conf
[Unit]
After=azure-complete-provisioning.service
Requires=azure-complete-provisioning.service
