# mkosi.images/buildernet-qemu/mkosi.extra/etc/systemd/system/persistent-setup.service.d/qemu-override.conf
[Service]
Environment=
ExecStart=
ExecStart=/usr/bin/persistent-setup --device-serial persistent --key-file /etc/disk-encryption-key
