# buildernet/mkosi.extra/etc/systemd/system/rbuilder-bidding-downloader.service
[Unit]
Description=Fetch rbuilder-bidding and top-bid-source binary
Wants=network-online.target
After=network-online.target render-config.service

[Service]
Type=oneshot
Environment=RBUILDER_BIDDING_SHA256SUM=208ebab8e6767b769136d25744805d906b6456f54084ca9f01e7a921db782dd5 \
            RBUILDER_BIDDING_TAG=v1.2.19 \
            TOP_BID_SOURCE_SHA256SUM=9e575caaf46110788665155e78b748a633bfc230fb147e25047b21a5d517f615 \
            TOP_BID_SOURCE_TAG=v1.2.19
EnvironmentFile=%E/rbuilder-bidding/gh_token.env
ExecStart=/bin/bash -c '[[ -f /usr/bin/rbuilder-bidding ]] || /usr/bin/fetch-gh-release-artifact fetch \
  --repo rbuilder-bidding-service --tag "$RBUILDER_BIDDING_TAG" \
  --artifact-name bidding-service --output-file /usr/bin/rbuilder-bidding'
ExecStart=/bin/bash -c '[[ -f /usr/bin/top-bid-source ]] || /usr/bin/fetch-gh-release-artifact fetch \
  --repo rbuilder-bidding-service --tag "$TOP_BID_SOURCE_TAG" \
  --artifact-name top-bid-source --output-file /usr/bin/top-bid-source'
ExecStart=/bin/bash -c '[[ -f /usr/bin/rbuilder-bidding ]] || \
  echo "$RBUILDER_BIDDING_SHA256SUM" /usr/bin/rbuilder-bidding | sha256sum --check'
ExecStart=/bin/bash -c '[[ -f /usr/bin/top-bid-source ]] || \
  echo "$TOP_BID_SOURCE_SHA256SUM" /usr/bin/top-bid-source | sha256sum --check'
ExecStart=/bin/chmod +x /usr/bin/rbuilder-bidding /usr/bin/top-bid-source
RemainAfterExit=yes
