# mkosi.images/gcp/mkosi.extra/etc/systemd/system/cvm-reverse-proxy-client.service.d/gcp-override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/cvm-reverse-proxy-client \
  --listen-addr=localhost:7937 \
  --target-addr=${BUILDERNET_BUILDERHUB_URL} \
  --client-attestation-type=dcap-tdx
SupplementaryGroups=
ProtectSystem=strict
ProtectHome=yes
AmbientCapabilities=CAP_DAC_OVERRIDE
ReadWritePaths=/sys/kernel/config/tsm/report
