# buildernet/mkosi.extra/etc/systemd/system/systemd-repart.service.d/90-persistent.conf
[Service]
Environment=BUILDERNET_TARGET_LUN=__BUILDERNET_TARGET_LUN
ExecStart=
ExecStart=/usr/bin/persistent-setup $BUILDERNET_TARGET_LUN
