# mkosi.profiles/playground/mkosi.extra/etc/systemd/system/operator-api.service.d/playground-override.conf
[Service]
ExecStartPre=+/usr/bin/sed -i 's/tls_create_if_missing = false/tls_create_if_missing = true/' /etc/operator-api/config.toml