# modules/l2/_common/mkosi.extra/etc/systemd/system/systemd-networkd.service.d/override.conf
[Service]
Environment=SYSTEMD_LOG_LEVEL=debug
Environment=SYSTEMD_LOG_LOCATION=1
