# mkosi.profiles/devtools/mkosi.extra/etc/systemd/system/dropbear.service.d/override.conf
[Unit]
# In devtools, remove the hard firewall dependency so SSH always works.
# Production images keep Requires=searcher-firewall.service via bob-common.
After=wait-for-key.service
Requires=wait-for-key.service
