# buildernet/mkosi.extra/etc/systemd/system/ntpd-rs.service.d/90-override.conf
[Service]
Restart=on-failure
AmbientCapabilities=
AmbientCapabilities=CAP_SYS_TIME
User=
Group=
DynamicUser=yes
