# modules/flashbox/observability/mkosi.extra/etc/systemd/system/searcher-firewall.service.d/needs-observability.conf
[Unit]
After=flashbox-observability-setup.service
Wants=flashbox-observability-setup.service
