# modules/flashbox/observability/mkosi.extra/etc/systemd/system/searcher-firewall.service.d/needs-observability.conf
[Unit]
After=fetch-observability-config.service
Wants=fetch-observability-config.service
